# etc/supervisord/ssr.conf
[supervisord]
nodaemon=true

[program:ssr]
directory=/var/www/html
command=nodemon --watch build --exec "node dist/server/server.js"
